# service/avacuationd.service
[Unit]
Description=Automatic Evacuation daemon

[Service]
ExecStart=/usr/bin/python /usr/share/evac/server.py