# service/evacuationd.service
[Unit]
Description=Automatic Evacuation daemon

[Service]
ExecStart=/usr/bin/python /usr/lib/python2.7/site-packages/evacuationd/server.py